# profile.d/000-locale-setup.sh
#!/bin/sh

# So Flatpak apps are installed with English and Spanish dictionaries
export LANG=en_US.UTF-8
export LANGUAGE=en_US.UTF-8:es_MX.UTF-8
