# profile.d/013-path-vscode.sh
#!/bin/sh

if [ -d "/Applications/Visual Studio Code.app/Contents/Resources/app/bin" ]; then
    pathprepend "/Applications/Visual Studio Code.app/Contents/Resources/app/bin"
fi
